# Backend/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'

[program:dotnet]
command=dotnet Tictactoe.dll
directory=/app